# deploy/systemd/frontend-interview-bank-admin.service.example
[Unit]
Description=Frontend Interview Bank Admin Static Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/frontend-interview-bank/current
ExecStart=/usr/local/bin/node scripts/serve-static.mjs --root /var/www/frontend-interview-bank/current/admin-web/dist --port 4174
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
